# root/etc/s6-overlay/s6-rc.d/sonarr-metadata-rewrite/run
#!/usr/bin/with-contenv bash

# Source environment
source /etc/environment

# Set default values if not provided
REWRITE_ROOT_DIR=${REWRITE_ROOT_DIR:-/tv}
PREFERRED_LANGUAGES=${PREFERRED_LANGUAGES:-en-US}
CACHE_DURATION_HOURS=${CACHE_DURATION_HOURS:-720}
PERIODIC_SCAN_INTERVAL_SECONDS=${PERIODIC_SCAN_INTERVAL_SECONDS:-86400}

# Export environment variables for the application
export REWRITE_ROOT_DIR
export PREFERRED_LANGUAGES
export CACHE_DURATION_HOURS
export PERIODIC_SCAN_INTERVAL_SECONDS
export CACHE_DIR=/config/cache
export ORIGINAL_FILES_BACKUP_DIR=/config/backups

# Change to app directory and run the application
cd /app
exec s6-setuidgid abc /app/.venv/bin/sonarr-metadata-rewrite